# rootfs/etc/s6-overlay/s6-rc.d/acars_stats/run
#!/bin/sh
exec /etc/s6-overlay/scripts/acars_stats
