# rootfs/etc/s6-overlay/s6-rc.d/irdm_stats/run
#!/bin/sh
exec /etc/s6-overlay/scripts/irdm_stats.sh
