# rootfs/etc/s6-overlay/s6-rc.d/hfdl_server/run
#!/bin/sh
exec /etc/s6-overlay/scripts/hfdl_server.sh
