# rootfs/etc/services.d/logsort/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

/scripts/line_trim.sh

sleep 60