# rootfs/etc/services.d/logsort/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

if [ -n "${VERBOSE}" ]; then
	/scripts/line_trim.sh
	sleep 60
else
    sleep 86400
fi