# rootfs/etc/s6-overlay/s6-rc.d/vdlm2_server/run
#!/bin/sh
exec /etc/s6-overlay/scripts/vdlm2_server.sh
